# vlasia_blog.service
[Unit]
Description=Vlasia Blog Gunicorn daemon
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/app  # Replace with actual path
Environment="PATH=/path/to/your/app/venv/bin"  # Replace with actual path
ExecStart=/path/to/your/app/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 wsgi:app
ExecReload=bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
